# demo/demo-init.conf
# Place in /etc/init

description "Demo"

start on net-device-up

respawn

exec java -jar /opt/demo.jar


